# packages/luci-app-arx-software/root/etc/uci-defaults/99-arx-opkg-mirror
#!/bin/sh
[ -x /sbin/arx-opkg-mirror ] && /sbin/arx-opkg-mirror default-firstboot
exit 0
