# packages/luci-app-arx-flash/root/etc/uci-defaults/99-arx-flash
#!/bin/sh
uci -q get arx-flash.main.enabled >/dev/null 2>&1 || {
	uci set arx-flash.main=main 2>/dev/null
	uci set arx-flash.main.enabled='1' 2>/dev/null
	uci commit arx-flash 2>/dev/null
}
exit 0
